# boc-warm.timer
[Unit]
Description=Daily BOC exchange rate cache warm

[Timer]
OnCalendar=*-*-* 10:40:00
Persistent=true

[Install]
WantedBy=timers.target
