# boc-warm.timer
[Unit]
Description=Periodic BOC exchange rate cache warm

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
